# scripts/systemd/overwatch-telegram.service
[Unit]
Description=Overwatch Telegram Bot
After=network.target

[Service]
Type=simple
User=overwatch
WorkingDirectory=/opt/overwatch
ExecStart=/usr/bin/node dist/telegram/index.js
Restart=always
RestartSec=5
EnvironmentFile=/opt/overwatch/.env

[Install]
WantedBy=multi-user.target
